# rootfs_skeleton/var/service/metalog/run
#! /bin/sh

exec 2>&1
exec metalog -N -C /etc/metalog.conf
