# rootfs_skeleton/etc/inittab
::sysinit:/bin/start.sh
::respawn:/bin/telnetd -l /bin/login -f /etc/banner.txt
