# rootfs_skeleton/etc/inittab
::sysinit:/bin/start_net.sh
::respawn:/bin/telnetd -l /bin/login -f /etc/banner.txt -F
