# rootfs_skeleton/var/service/lighttpd/run
#!/bin/sh

exec 2>&1
exec lighttpd -D -f /etc/lighttpd.conf
